# rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

mkdir -p /tmp/nginx
mkdir -p /app/logs
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf

