# rootfs/etc/services.d/geminabox/run
#!/usr/bin/with-contenv sh

cd /app
bundle exec unicorn -c lib/unicorn.rb

